# subber.conf
# subber.conf
start on filesystem
exec /home/hakon/Subber/start.sh