# subber.conf
# subber.conf
start on filesystem
task
exec /home/hakon/Subber/start.sh