# subber.conf
description     "Start subber api server"
# subber.conf
start on filesystem
respawn
task
exec /home/hakon/Subber/start.sh >> /home/hakon/Subber/subber.log
